# alerting.rules
ACTION=="change", SUBSYSTEM=="power_supply", KERNEL=="BAT1", RUN+="/home/ink/Documents/GitHub/alerting/script.sh"
