# systemd/alerting.service
[Unit]
Description=Run alerting script

[Service]
Type=oneshot
ExecStart=/home/ink/Documents/GitHub/alerting/reference.sh
